# supervisor/conf.d/datastax-agent.conf
[program:datastax-agent]
command=/opt/datastax-agent/bin/datastax-agent -f
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

